# systemd/trackora.service
[Unit]
Description=Trackora background tracking service
Documentation=file:///usr/share/doc/trackora/systemd-service.md
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m trackora
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillMode=process
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=default.target
